# debian/postrm
#!/bin/bash

systemctl disable gooroom-shutdown-trigger

[ -e /etc/systemd/system/gooroom-shutdown-trigger.service ] && rm /etc/systemd/system/gooroom-shutdown-trigger.service

